# jobs/health.service
[Unit]
Description=Check for MCE hardware errors

[Service]
Type=oneshot
ExecStart=/home/felipeftn/dotfiles/jobs/health.sh
